# config/webconfig.service
[Unit]
Description="Web config for pharma"
Wants=network.target
After=network.target


[Service]
Type=simple
GuessMainPID=yes
ExecStart=/bin/bash -c 'python webctrl.py'
ExecStop=kill -SIGINT $mainpid
Restart=on-failure
WorkingDirectory=/home/compost/pharma-main-unit/src/
# affichage des log en temps reel dans journactl
Environment=PYTHONUNBUFFERED=1

[Install]
Alias=PharmaWebConfig
WantedBy=multi-user.target


